# debian_package/git2consul.service
[Unit]
Description=git2consul
After=syslog.target
After=network.target
ConditionPathExists=/usr/share/git2consul

[Service]
Type=simple
User=git2consul
Group=git2consul
ExecStartPre=/usr/share/git2consul/utils/config_seeder.js /etc/git2consul/config.json
ExecStart=/usr/bin/nodejs /usr/share/git2consul/lib/index.js
TimeoutStopSec=59
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=git2consul

[Install]
WantedBy=multi-user.target
